# deploy/openrc
#!/sbin/openrc-run

name="Abaku Capital website"
description="Static website server for Abaku Capital"

# OpenRC may set $0 to its internal runner, so resolve the installed service
# symlink explicitly instead.
service_file="$(readlink -f "/etc/init.d/$RC_SVCNAME")"
default_app_dir="$(dirname "$(dirname "$service_file")")"

: "${ABAKU_DIR:=$default_app_dir}"
: "${ABAKU_USER:=www-data}"
: "${ABAKU_GROUP:=www-data}"
: "${ABAKU_NODE:=/usr/bin/node}"
: "${ABAKU_LOG:=/var/log/abaku.capital.log}"
: "${ABAKU_PORT:=4173}"

export PORT="$ABAKU_PORT"

directory="$ABAKU_DIR"
command="$ABAKU_NODE"
command_args="server.js"
command_user="$ABAKU_USER:$ABAKU_GROUP"

supervisor="supervise-daemon"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="$ABAKU_LOG"
error_log="$ABAKU_LOG"
respawn_delay=5
respawn_max=0

depend() {
  need net
  after firewall
}

start_pre() {
  case "$ABAKU_PORT" in
    ''|*[!0-9]*)
      eerror "ABAKU_PORT must be a number between 1 and 65535"
      return 1
      ;;
  esac

  if [ "$ABAKU_PORT" -lt 1 ] || [ "$ABAKU_PORT" -gt 65535 ]; then
    eerror "ABAKU_PORT must be a number between 1 and 65535"
    return 1
  fi

  checkpath --file --mode 0644 --owner "$command_user" "$ABAKU_LOG"

  if [ ! -r "$ABAKU_DIR/server.js" ]; then
    eerror "Cannot read $ABAKU_DIR/server.js"
    return 1
  fi

  if [ ! -x "$ABAKU_NODE" ]; then
    eerror "Node executable not found at $ABAKU_NODE"
    return 1
  fi
}
